# shell/com.giwooklee.trading-system.daily-news-agy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.giwooklee.trading-system.daily-news-agy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/giwooklee/Workspace/trading_system/shell/run_daily_news_summary_agy.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/giwooklee/Workspace/trading_system/logs/daily_news_agy.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/giwooklee/Workspace/trading_system/logs/daily_news_agy.launchd.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/giwooklee/Workspace/trading_system</string>
</dict>
</plist>
